# proxy/supervisor/conf.d/free_proxy_list_net.conf
[program:free_proxy_list_net]
directory=/app
command=python3 /app/free_proxy_list_net.py
autostart=true
autorestart=true
redirect_stderr=true
