# proxy/supervisor/conf.d/ip_address_com.conf
[program:ip_address_com]
directory=/app
command=python3 /app/ip_address_com.py
autostart=true
autorestart=true
redirect_stderr=true
