# proxy/supervisor/conf.d/verify_proxy.conf
[program:verify_proxy]
directory=/app
command=python3 /app/verify_proxy.py
autostart=true
autorestart=true
redirect_stderr=true
