# rds/supervisor/conf.d/page.conf
[program:page]
directory=/app
command=python3 /app/page.py
autostart=true
autorestart=true
redirect_stderr=true
