# rds/supervisor/conf.d/pagination.conf
[program:pagination]
directory=/app
command=python3 /app/pagination.py
autostart=true
autorestart=true
redirect_stderr=true
